# .config/systemd/user/games-backup.timer
[Unit]
Description=Run game save backup every 10 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
