# .config/systemd/user/github-backup@.service
[Unit]
Description=Automated Git Backup for %i

[Service]
Type=oneshot
ExecStart=/home/aechewhy/.local/bin/github-backup.sh %i
