# .config/systemd/user/github-backup@.timer
[Unit]
Description=10-minute backup timer for %i

[Timer]
OnBootSec=5m
OnUnitActiveSec=10m
Persistent=true

[Install]
WantedBy=timers.target
