# .config/systemd/user/water-reminder.timer
[Unit]
Description=Hydration Reminder Timer

[Timer]
OnActiveSec=5m
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
