# .config/systemd/user/water-reminder.service
[Unit]
Description=Hydration Notification Trigger

[Service]
Type=oneshot
ExecStart=%h/.local/bin/water-reminder.sh
